# roles/nut-rpt1000ap/templates/90-nut-ups.rules.j2
# {{ ansible_managed }}

SUBSYSTEM=="usb", ACTION=="add", ATTR{idVendor}=="0d9f", ATTR{idProduct}=="0004", TAG+="systemd", RUN+="/usr/bin/systemctl restart nut-server"
